# ObiOra-Suite/crashhunter/systemd/crashhunter.service
[Unit]
Description=Crash Hunter — Linux Freeze Diagnosis Daemon (Black Box Flight Recorder)
Documentation=https://obiora.io
After=network-online.target libvirtd.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
Environment=CRASHHUNTER_BASE=/opt/crashhunter
Environment=CRASHHUNTER_INTERVAL=5
Environment=CRASHHUNTER_RING_CAPACITY=720
Environment=CRASHHUNTER_LOG_LEVEL=INFO
ExecStart=/opt/crashhunter/venv/bin/crashhunter run
Restart=always
RestartSec=10
# Allow graceful ftrace/subprocess cleanup on systemctl stop (SIGTERM → cleanup → exit).
TimeoutStopSec=20
Nice=10
IOSchedulingClass=idle
MemoryMax=128M
CPUQuota=10%

# Hardening
NoNewPrivileges=false
ProtectSystem=strict
ReadWritePaths=/opt/crashhunter /var/log /proc /sys /run
ProtectHome=true

[Install]
WantedBy=multi-user.target
